# bin/osm.service
# systemd unit file
#
# place in /etc/systemd/system
# systemctl enable snt.service
# systemctl start snt.service

[Unit]
Description=Service Tomorrow Server
After=network.target
[Service]
WorkingDirectory=/var/www/osm/packages/core
ExecStart=/home/aj/.nvm/versions/node/v11.13.0/bin/node /var/www/osm/packages/core/dist
[Install]
WantedBy=multi-user.target
